# deploy/super-agent.service.example
[Unit]
Description=Super Agent one-shot production runner
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=super-agent
Group=super-agent
WorkingDirectory=/var/lib/super-agent
EnvironmentFile=/etc/super-agent/super-agent.env
ExecStart=/usr/bin/node /opt/super-agent/scripts/linux-sandbox-launcher.mjs -- /usr/bin/node /opt/super-agent/dist/bin/super-agent.js run --session production-job --prompt "Replace with an operator supplied prompt"

# Required by the runtime release contract. The launcher moves itself into a
# runtime leaf, leaving this delegated service root process-free for operation
# cgroups. On launcher/Agent death systemd kills every remaining descendant.
Delegate=cpu memory pids
KillMode=control-group
# Make the final descendant sweep explicit instead of relying on systemd defaults.
SendSIGKILL=yes
OOMPolicy=stop
LimitNOFILE=4096
UMask=0077

# The parent service envelope must be larger than any one operation cgroup.
MemoryMax=2G
MemorySwapMax=0
TasksMax=256
CPUQuota=200%

Restart=no
TimeoutStopSec=10s

[Install]
WantedBy=multi-user.target
